# deploy/edgevest-poller.service
[Unit]
Description=Edgevest Live Poller
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/edgevest
EnvironmentFile=/home/ubuntu/edgevest/.env
ExecStart=/home/ubuntu/edgevest/edgevest-poller
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
